# roles/landoop.topics-ui/templates/topics-ui.service
[Unit]
Description=Landoop Kafka Topics UI

[Service]
Type=notify
ExecStart=http-server -p 8080 .

[Install]
WantedBy=multi-user.target
